# scripts/com.claude.bob-healthcheck.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.bob-healthcheck</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/andykaufman/.claude/scripts/bob-healthcheck.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/andykaufman</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/andykaufman/.claude/scheduled-tasks/bob-healthcheck/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/andykaufman/.claude/scheduled-tasks/bob-healthcheck/launchd-stderr.log</string>
    <key>RunAtLoad</key>
    <false/>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
